# packaging/systemd/mai-api.service
[Unit]
Description=Island Mountain MAI API server (REST + gRPC)
Documentation=file:///usr/share/doc/mai/README.md
Documentation=https://github.com/USS-Parks/Mighty-Eel-OS
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/var/lib/mai /etc/mai

[Service]
Type=notify
User=mai
Group=mai
WorkingDirectory=/var/lib/mai
Environment=MAI_SHIP_PROFILE=/etc/mai/profile.toml
Environment=RUST_LOG=mai_api=info,mai_core=info,mai_hil=info,warn
ExecStartPre=/usr/bin/mai-ship-validate --profile /etc/mai/profile.toml
ExecStart=/usr/bin/mai-api --config /etc/mai/profile.toml
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60s
TimeoutStartSec=120s
TimeoutStopSec=30s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mai-api

# Privilege hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ReadWritePaths=/var/lib/mai /var/log/mai /run/mai /var/backups/mai
ReadOnlyPaths=/etc/mai

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
TasksMax=infinity

[Install]
WantedBy=multi-user.target
